# docker/benchmarks/benchmarks-server.service
[Unit]
Description=ASP.NET Benchmarks Server
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
WorkingDirectory=/opt/benchmarks-server/docker/benchmarks
Environment=hardware_version=HPZ440
ExecStartPre=/usr/bin/git -C ../../ pull origin master
ExecStartPre=/bin/bash -xe build.sh
ExecStart=/bin/bash -xe run-server-in-systemd.sh

[Install]
WantedBy=multi-user.target
